# resources/minecraft.socket
[Unit]
Description=minecraft reads commands from this

[Socket]
SocketUser=ec2-user
Service=minecraft.service
ListenFIFO=/tmp/minecraft.stdin
